# _SPIdumpTPlink.bin.extracted/squashfs-root/etc/hotplug.d/net/20-client_mgmt
#!/bin/sh

[ "$TYPE" = "port" -a "$ACTION" = "change" -o "$TYPE" = "wsta" ] && {
	ubus call client_mgmt update '{"delay":3000}'
}