# _SPIdumpTPlink.bin.extracted/squashfs-root/etc/hotplug.d/net/21-switch_ports
#!/bin/sh

[ "$TYPE" = "port" -a "$ACTION" = "change" ] && {
	. /usr/lib/switch_led/switch_led.sh
	notify_portlink_changed
	
	notify_online_test
}