# _SPIdumpTPlink.bin.extracted/squashfs-root/etc/hotplug.d/firewall/65-smart_home
#!/bin/sh

. /lib/functions.sh

[ "$ACTION" = "start" ] && {
    /etc/init.d/smart_home restart
}

