# _SPIdumpTPlink.bin.extracted/squashfs-root/etc/hotplug.d/iface/15-uhttpd
#!/bin/sh

if [ "$INTERFACE" == "lan" -o "$INTERFACE" == "lanv6" ] && [ "$ACTION" == "ifup" ]; then
	echo "======== uhttpd: lan or lanv6 update ========"
	ubus call uhttpd lan_update '{}'
fi
