# _SPIdumpTPlink.bin.extracted/squashfs-root/etc/hotplug.d/iface/99-smart_home
#!/bin/sh

. /lib/functions.sh

[ "$ACTION" == "ifup" ] && [ "$INTERFACE" == "wan" -o "$INTERFACE" == "internet" ] && {
    /etc/init.d/smart_home restart
}

